# modules/ssh-server/home/.config/dotfiles/profile.d/60-ssh-server.sh
# Server-profile shell additions.
export DOTFILES_SERVER_PROFILE=1

dotfiles_source_optional "$XDG_CONFIG_HOME/dotfiles/local.server.sh"
